# crontab
* * * * * echo "Hello world" >> /dev/stdout 2>&1
